# lab3/services/notes_app.service
[Unit]
Description=Notes App
After=network.target

[Service]
User=mikhail
WorkingDirectory=/home/mikhail/DevOps/lab3/notes_app
ExecStart=/usr/bin/python3 app.py
Restart=always

[Install]
WantedBy=multi-user.target
